# package/debian/job-run/DEBIAN/postinst
#!/bin/bash

echo "creating folders"

mkdir -p /var/log/urmia
chown -R urmia:urmia /var/log/urmia

mkdir -p /var/run/urmia
chown -R urmia:urmia /var/run/urmia

echo "installing job run service"
/bin/systemctl enable urmia-job-run.service

exit 0

