# package/debian/obj-store/DEBIAN/prerm
#!/bin/bash

echo "disabling service"
/usr/sbin/service urmia-storage status
[ $? -eq 0 ] && /usr/sbin/service urmia-storage stop

/bin/systemctl disable urmia-storage.service

exit 0