# package/debian/obj-api/DEBIAN/prerm
#!/bin/bash

echo "disabling obj api service"
/usr/sbin/service urmia-obj-api status
[ $? -eq 0 ] && /usr/sbin/service urmia-obj-api stop

/bin/systemctl disable urmia-obj-api.service

exit 0